# crontab
# m h  dom mon dow   command
00 10 * * 4-7 "shutdown -r now"
